# deployment/radgnarackassist.service
[Unit]
Description=RadgnarackAssist Backend API
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/radgnarack-assist
Environment="PATH=/opt/radgnarack-assist/venv/bin"
EnvironmentFile=/opt/radgnarack-assist/.env

# Use port 8001 to avoid conflicts with other services
ExecStart=/opt/radgnarack-assist/venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8001

Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/radgnarack-assist

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=radgnarackassist

[Install]
WantedBy=multi-user.target
